# meta-axelera/recipes-security/optee/optee-rockchip/tee-supplicant.service
[Unit]
Description=TEE Supplicant

[Service]
User=root
ExecStart=/usr/bin/tee-supplicant

[Install]
WantedBy=basic.target
